# upstart_dvdpila.conf
# DVD Pila! - DVD Archive Server
#
# Organize your DVD archives.

description	"DVD Pila! server"

start on runlevel [2345]
stop on runlevel [!2345]

exec /home/adam/work/dvdpila/tornado_server.py

